# service_files/weakbackup.service
[Unit]
Description=Weakbackup
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker start WB-restic
ExecStop=/usr/bin/docker stop WB-restic

[Install]
WantedBy=local.target