# service_files/weakbackup.timer
[Unit]
Description=Weakbackup Timer
Requires=docker.service

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true